# assets/hefesto.service
[Unit]
Description=Hefesto DualSense adaptive trigger daemon
After=graphical-session.target
PartOf=graphical-session.target
# BUG-MULTI-INSTANCE-01: teto anti-loop em Unit. Máximo 3 tentativas em 30s —
# se estourar, systemd marca failed e pára de respawnar.
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
ExecStart=%h/.local/bin/hefesto daemon start --foreground
Restart=on-failure
RestartSec=2
# BUG-MULTI-INSTANCE-01: SIGTERM (143) é desligamento limpo — o takeover da
# próxima instância e o "Sair" do tray não devem disparar respawn.
SuccessExitStatus=143 SIGTERM
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=graphical-session.target
