# assets/hefesto-headless.service
[Unit]
Description=Hefesto DualSense daemon (headless — sem auto-switch X11)
After=default.target
Conflicts=hefesto.service

[Service]
Type=simple
ExecStart=%h/.local/bin/hefesto daemon start --foreground --headless
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1
Environment=HEFESTO_NO_WINDOW_DETECT=1

[Install]
WantedBy=default.target
