# assets/hefesto-dualsense4unix.service
[Unit]
Description=Hefesto - Dualsense4Unix DualSense adaptive trigger daemon
# BUG-DEB-AUTOSTART-WANTEDBY-DEFAULT-01: o daemon usa /dev/hidraw e evdev,
# não depende de DISPLAY/Wayland/X11. Ancorar em default.target garante
# autostart resiliente em qualquer DE (GNOME/KDE/Sway/i3) e elimina race
# entre login e graphical-session activation. After= preservado por ordem.
After=graphical-session.target default.target
# BUG-MULTI-INSTANCE-01: teto anti-loop em Unit. Máximo 3 tentativas em 30s —
# se estourar, systemd marca failed e pára de respawnar.
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
ExecStart=%h/.local/bin/hefesto-dualsense4unix daemon start --foreground
Restart=on-failure
RestartSec=2
# BUG-MULTI-INSTANCE-01: SIGTERM (143) é desligamento limpo — o takeover da
# próxima instância e o "Sair" do tray não devem disparar respawn.
SuccessExitStatus=143 SIGTERM
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
