# assets/hefesto-dualsense4unix-gui-hotplug.service
[Unit]
Description=Abre a GUI do Hefesto - Dualsense4Unix ao conectar o DualSense (hotplug USB)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
# BUG-TRAY-SINGLE-FLASH-01: guard pgrep removido — era race-prone com udev
# ADD duplo em <200ms. O próprio launcher usa acquire_or_bring_to_front:
# se já há GUI ativa, traz-a ao foco e sai com exit 0.
# %h expande para o HOME do usuário dono da sessão systemd --user.
ExecStart=%h/.local/bin/hefesto-dualsense4unix-gui

[Install]
WantedBy=graphical-session.target
